# puppet/files/etc/init/statsd.conf
# This Source Code Form is subject to the terms of the Mozilla Public
# License, v. 2.0. If a copy of the MPL was not distributed with this
# file, You can obtain one at http://mozilla.org/MPL/2.0/.

# cat /etc/init/statsd.conf
description "statsd"
author      "lorchard"

start on startup
stop on shutdown

expect daemon
respawn

script  
    export HOME="/home/vagrant"
    exec sudo -u vagrant node /home/vagrant/node_modules/.bin/statsd /home/vagrant/statsd-config.js
end script
