# puppet/files/etc/init/statsd.conf
# cat /etc/init/statsd.conf
description "statsd"
author      "lorchard"

start on startup
stop on shutdown

expect daemon
respawn

script  
    export HOME="/home/vagrant"
    exec sudo -u vagrant node /home/vagrant/node_modules/.bin/statsd /home/vagrant/statsd-config.js
end script
